# install/decentered_import.timer
[Unit]
Description=Import events into the decentered_events db

[Timer]
OnCalendar=*-*-* *:30
Persistent=true

[Install]
WantedBy=timers.target